# systemd/notes-api.service
[Unit]
Description=Notes Pipeline API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/notes-pipeline
EnvironmentFile=/opt/notes-pipeline/.env
ExecStart=/opt/notes-pipeline/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
